# myapp_gunicorn_psql_deb/etc/systemd/system/crl_daemon.service
[Unit]
Description=CRL Daemon Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/myapp/crl_daemon
Environment="PATH=/opt/myapp/venv/bin"
ExecStart=/opt/myapp/venv/bin/python3 /opt/myapp/crl_daemon/crl_daemon.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=crl_daemon

[Install]
WantedBy=multi-user.target